# private_dot_config/autostart/autorandr.desktop
[Desktop Entry]
Type=Application
Name=autorandr
Exec=bash -c "sleep 3 && autorandr --change"
X-GNOME-Autostart-Phase=Applications
NoDisplay=true
