# private_dot_config/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=Network (nm-applet)
Exec=/usr/bin/nm-applet
Hidden=true
X-GNOME-Autostart-enabled=false
